# traefik/com.traefik.startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.traefik.startup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>source /private/etc/traefik/.env && exec /usr/local/bin/traefik</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/traefik.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/traefik.error.log</string>
    <key>WorkingDirectory</key>
    <string>/private/etc/traefik</string>
</dict>
</plist> 